# chest/crontab
MAILTO=""
SHELL=/bin/bash
*/10 * * * * /Users/caio.tavares/dev/functions/chest_job.sh &> /Users/caio.tavares/dev/functions/chest_job.log
