# chest/crontab
MAILTO=""
SHELL=/bin/bash
*/10 * * * * /Users/caio.tavares/dev/personal/chest/chest/job.sh &> /Users/caio.tavares/dev/personal/chest/chest/job.log
